# cloudflare-exporter.timer
[Unit]
Description=Cloudflare Exporter Timer

[Timer]
OnCalendar=*-*-* 23:55:00

[Install]
WantedBy=timers.target